# examples/supervisor.conf
[program:gumicorn]
command=/usr/local/bin/gumicorn main:application -c /path/to/project/gumicorn.conf.py
directory=/path/to/project
user=nobody
autorestart=true
redirect_stderr=true
